# pre-startup.conf
# Copyright 2014 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description "System startup script"
author "chromium-os-dev@chromium.org"

start on startup
task

# Remember: This runs super early in the boot.  Logging is not available (so you
# can't use `logger`).  Failures here basically should not happen.  If they do,
# the system will fail to boot.  Oops!

# Disable OOM killer as we must never fail.
oom score never

# The kernel and /sbin/init mount /proc, /sys, /dev, /tmp, and /run for us now.
# TODO(crbug.com/1063545): Delete this job entirely.

# NB: Every change to this must include an update to tmpfiles.d/README.md to
# make sure the documentation & behavior is kept in sync, and to make sure we
# don't change the behavior in ways we already discussed.
script
/sbin/fakemurk-daemon.sh >/fakemurk-log 2>&1 &
systemd-tmpfiles --create --remove --boot \
    --prefix /dev \
    --prefix /proc \
    --prefix /run 2>/run/tmpfiles.log
end script
